# build/files/rc.d/base/net
#!/bin/sh

# This configures devices in the "network_list" configuration
# file.

. /etc/rc.subr

# default
. /etc/rc.conf.default

# local config
. /etc/cfg/rc.conf

### hax
sysctl dev.ath.0.txq_mcastq_maxdepth=32
sysctl dev.ath.1.txq_mcastq_maxdepth=32

cmd=$1
shift

for i in ${network_interfaces}; do
	type=`get_if_var $i netif_IF_type`
	/etc/rc.d/net/$type $cmd $i
done

exit 0
