# build/files/rc.d/net/vlan
#!/bin/sh

# This configures a vlan subinterface.

. /etc/rc.subr
. /etc/rc.conf.default
. /etc/cfg/rc.conf

# variables:
#
# netif_IF_type="vlan"
# netif_IF_enable="yes|no"
# netif_IF_addrtype="static|dhcp|none"
# netif_IF_descr="description"
# netif_IF_name="arge0.2"
# netif_IF_parent="arge0"
# netif_IF_vlanid="2"
#
# for static addresses:
# netif_IF_ipv4_address="x.x.x.x"
# netif_IF_ipv4_netmask="y.y.y.y"

cmd=$1
shift

ifname=$1
shift

interface_enable=`get_if_var $ifname netif_IF_enable no`
interface_addrtype=`get_if_var $ifname netif_IF_addrtype`
interface_descr=`get_if_var $ifname netif_IF_descr`
interface_name=`get_if_var $ifname netif_IF_name`
interface_parent=`get_if_var $ifname netif_IF_parent`
interface_vlanid=`get_if_var $ifname netif_IF_vlanid`

interface_ipv4_address=`get_if_var $ifname netif_IF_ipv4_address`
interface_ipv4_netmask=`get_if_var $ifname netif_IF_ipv4_netmask`

case $cmd in
	'start')
		echo "*** Interface: ${ifname}: start"
		ifconfig ${interface_name} create vlan ${interface_vlanid} vlandev ${interface_parent} || exit 1
		ifconfig ${interface_name} up
		if [ "x${interface_addrtype}" = "xstatic" ]; then
			ifconfig ${interface_name} inet ${interface_ipv4_address} \
			    netmask ${interface_ipv4_netmask}
		elif [ "x${interface_addrtype}" = "xdhcp" ]; then
			dhclient ${interface_name} &
		fi
		echo "*** Interface: ${ifname}: done"
		exit 0
	;;
	'stop')
		echo "*** Interface: ${ifname}: stop"
		ifconfig ${interface_name} down
		ifconfig ${interface_name} destroy
		echo "*** Interface: ${ifname}: done"
		exit 0
	;;
	*)
		echo "Usage: $0 {start|stop} {name}"
esac

exit 0
